# supervisord.conf/homesite.conf
[program:homesite]
command = sh -c 'sleep 5 && exec node ./bin/www'
autostart = true
stdout_logfile = /homesite/logs/homesite.log
redirect_stderr = true
autorestart = true
